# cookbooks/bcpc/templates/default/calico/calico-felix.service
[Unit]
Description=Calico Felix agent
After=syslog.target network.target

[Service]
User=root
Environment=ETCD_ENDPOINTS=http://localhost:2379
ExecStartPre=/bin/mkdir -p /var/run/calico
ExecStart=/usr/bin/calico-felix
KillMode=process
Restart=on-failure
LimitNOFILE=32000

[Install]
WantedBy=multi-user.target
