# cookbooks/bcpc-consul/templates/default/consul.service.erb
[Unit]
Description=Consul daemon
Requires=network-online.target

[Service]
Type=simple
EnvironmentFile=~/etc/default/consul
User=<%= @username %>
ExecStart=<%= @executable %> agent $OPTIONS -config-dir=<%= @conf_dir %>
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal
SyslogIdentifier=consul

[Install]
WantedBy=multi-user.target
